# common/backup/.config/systemd/user/backup.service
[Unit]
Description=backup org files

[Service]
Type=oneshot
ExecStart=%h/.local/bin/vsimoes-backup.sh
